# services/greek-tutor-fastapi.service
[Unit]
Description=Greek Tutor FastAPI (Uvicorn)
After=network.target

[Service]
Type=simple
# Edit these paths before installing:
WorkingDirectory=/opt/greek_tutor
EnvironmentFile=/etc/default/greek-tutor
ExecStart=/bin/bash -lc '"${VENV_DIR}/bin/uvicorn" fastapi_server.api:app --host ${API_HOST:-127.0.0.1} --port ${API_PORT:-8000} --workers ${API_WORKERS:-2}'
Restart=on-failure
RestartSec=3
# Optionally set a runtime user:
# User=greek

[Install]
WantedBy=multi-user.target

