# services/greek-tutor-fastapi.service
[Unit]
Description=Greek Tutor FastAPI (Uvicorn)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/zacharyclement/greek_tutor
EnvironmentFile=/etc/default/greek-tutor
ExecStart=/bin/bash -c '/home/zacharyclement/greek_tutor/.venv/bin/uvicorn fastapi_server.api:app --host $${API_HOST:-127.0.0.1} --port $${API_PORT:-8000} --workers $${API_WORKERS:-2}'
Restart=on-failure
RestartSec=3
User=zacharyclement

[Install]
WantedBy=multi-user.target

